# packaging/debian/postinst
#!/bin/sh
set -e

SERVICE="openhd-glide.service"
DEFAULTS="/etc/default/openhd-glide"
ETC_UNIT="/etc/systemd/system/$SERVICE"
OLD_SERVICE_RUN="/usr/local/lib/openhd-glide/openhd-glide-service-run"
SERVICE_RUN="/usr/lib/openhd-glide/openhd-glide-service-run"

is_enabled() {
  case "$(printf '%s' "$1" | tr '[:upper:]' '[:lower:]')" in
    1|true|yes|on|enabled)
      return 0
      ;;
    *)
      return 1
      ;;
  esac
}

if [ "$1" = "configure" ]; then
  if command -v systemctl >/dev/null 2>&1; then
    if [ -f "$ETC_UNIT" ] && [ -x "$SERVICE_RUN" ]; then
      sed -i "s#^ExecStart=$OLD_SERVICE_RUN\$#ExecStart=$SERVICE_RUN#" "$ETC_UNIT"
    fi

    if [ -d /run/systemd/system ]; then
      for obsolete_service in openhd-glide-autostart.service openhd-glide-stream.service openhd-video.service; do
        if systemctl list-unit-files "$obsolete_service" >/dev/null 2>&1; then
          systemctl disable --now "$obsolete_service" >/dev/null 2>&1 || true
        fi
      done
    fi

    systemctl daemon-reload >/dev/null 2>&1 || true

    if [ -d /run/systemd/system ]; then
      GLIDE_START_ON_INSTALL=0
      if [ -r "$DEFAULTS" ]; then
        set +e
        . "$DEFAULTS"
        set -e
      fi

      if is_enabled "${GLIDE_START_ON_INSTALL:-0}"; then
        systemctl restart "$SERVICE" >/dev/null 2>&1 || true
      elif systemctl is-active --quiet "$SERVICE"; then
        systemctl try-restart "$SERVICE" >/dev/null 2>&1 || true
      fi
    fi
  fi
fi

exit 0
