# packaging/systemd/openhd-glide.service.in
[Unit]
Description=OpenHD Glide KMS rendering stack
Documentation=https://github.com/OpenHD
After=systemd-udev-settle.service network.target
Wants=systemd-udev-settle.service
ConditionPathExists=/dev/dri
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
Type=simple
ExecStart=@OPENHD_GLIDE_SERVICE_RUN_FULL_PATH@
Restart=on-failure
RestartSec=2
KillMode=control-group
RuntimeDirectory=openhd-glide
RuntimeDirectoryMode=0755
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
